# debian/postinst
#!/bin/sh
set -e

# 获取实际安装用户（兼容sudo安装）
username=$(logname 2>/dev/null || echo "$SUDO_USER")
homedir=$(getent passwd "$username" | cut -d: -f6)

# 创建必要目录
mkdir -p "$homedir/.config/systemd/user"
mkdir -p "$homedir/.local/share/MemoCoco"

# 复制服务文件到用户目录
cp /usr/share/memococo/memococo.service "$homedir/.config/systemd/user/"

# 修复权限
chown -R "$username":"$username" "$homedir/.config/systemd/user/memococo.service"
chown -R "$username":"$username" "$homedir/.local/share/MemoCoco"

# 启用用户级服务
sudo -u "$username" systemctl --user daemon-reload
sudo -u "$username" systemctl --user enable memococo.service
sudo -u "$username" systemctl --user start memococo.service

# 允许用户未登录时保持运行
systemctl enable --now --user linger "$username"